# docs/Linux/crontab/cron log by day.crontab

#Envoi de mail AZ
30 9 * * * /usr/bin/php5-cgi /home/medica-source/www/az/cron-email-as.php > /var/log/az-mail-cron.`date +\%d-\%m-\%H-\%M`.log 2>&1
